# installation/services/openrc/immich-web
#!/sbin/openrc-run

: ${data_path:="/usr/src/web"}

name="Immich web"
directory="$data_path"
command=/bin/ash
command_args="-c $data_path/entrypoint.sh"

command_background=yes
command_user="node:node"
pidfile="/run/immich-web.pid"

output_log="/var/log/immich/web.log"
error_log="/var/log/immich/web.err.log"

required_dirs="$data_path"
required_files="$data_path/entrypoint.sh"

depend() 
{
        need net
}

# Load env variables
start_pre()
{
        source /etc/profile
}

# Find children process and kill them
stop_pre()
{
	if pgrep -P "$(cat $pidfile)" &> /dev/null
        then
                pkill -P "$(cat $pidfile)"
        fi
}
